# debian/snort-common.conffiles
/etc/cron.daily/5snort
/etc/snort/snort.conf
/etc/snort/unicode.map
/etc/snort/threshold.conf
