# debian/snort-mysql.conffiles
/etc/init.d/snort
/etc/ppp/ip-up.d/snort
/etc/ppp/ip-down.d/snort
/etc/logrotate.d/snort
/etc/default/snort
